# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x1024x24
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/xvfb.err.log
stdout_logfile=/var/log/supervisor/xvfb.out.log

[program:openbox]
command=/usr/bin/openbox --replace
autostart=true
autorestart=true
environment=DISPLAY=:1
stderr_logfile=/var/log/supervisor/openbox.err.log
stdout_logfile=/var/log/supervisor/openbox.out.log

[program:kasmvnc]
command=/opt/kasmvnc/bin/vncserver -fg
autostart=true
autorestart=true
environment=DISPLAY=:1
stderr_logfile=/var/log/supervisor/kasmvnc.err.log
stdout_logfile=/var/log/supervisor/kasmvnc.out.log

[program:pulseaudio]
command=/usr/bin/pulseaudio --start
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/pulseaudio.err.log
stdout_logfile=/var/log/supervisor/pulseaudio.out.log
